# deploy/super-scraper.service
[Unit]
Description=Super Scraper (Docker Compose)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=300
WorkingDirectory=/etc/super-scraper

ExecStartPre=/usr/bin/docker compose -f /etc/super-scraper/docker-compose.yml pull
ExecStart=/usr/bin/docker compose -f /etc/super-scraper/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /etc/super-scraper/docker-compose.yml down --remove-orphans

[Install]
WantedBy=multi-user.target

