# deploy/super-scraper-update.timer
[Unit]
Description=Check for Super Scraper image updates every 5 minutes

[Timer]
# First check 2 minutes after boot
OnBootSec=2min
# Then once daily at 4:00 AM
OnCalendar=*-*-* 04:00:00

[Install]
WantedBy=timers.target

